# homebox/rootfs/etc/s6-overlay/s6-rc.d/init-homebox/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant App: Homebox
# Configures Homebox before running
# ==============================================================================

bashio::log.info "Initializing Homebox..."

# Ensure data directory exists
mkdir -p /data

# Read and apply configuration options
# Example:
# if bashio::config.true 'ssl'; then
#     bashio::log.info "SSL is enabled"
# fi

bashio::log.info "Initialization complete"
